# technologies/ansible/day16-project-deployment/roles/monitoring/templates/alertmanager.service.j2
[Unit]
Description=Alertmanager
After=network-online.target

[Service]
User=alertmanager
Group=alertmanager
Type=simple
ExecStart=/usr/local/bin/alertmanager \
    --config.file={{ alertmanager_config_dir }}/alertmanager.yml \
    --storage.path={{ alertmanager_data_dir }} \
    --web.listen-address=:{{ alertmanager_port }}

[Install]
WantedBy=multi-user.target